# services/system-configurations/etc/systemd/system/docker.service.d/proxy.conf
[Service]
Environment="HTTP_PROXY=http://10.0.1.68:8889"
Environment="HTTPS_PROXY=http://10.0.1.68:8889"
Environment="NO_PROXY=localhost,127.0.0.1,nvcr.io,aliyuncs.com,edu.cn,cvgl.lab"